# target/linux/ramips/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
::respawn:/opt/ejoin_fs/script/mtk_autorun.sh
::askconsole:/bin/ash --login
